# salt/states/common/beaver/files/etc/init/beaver.conf
description "ship logs to logstash via beaver"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec /usr/bin/beaver -c /etc/beaver/beaver.conf -C /etc/beaver/conf.d -t redis
